# deploy-setting/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=service php7.0-fpm start
autostart=true
autorestart=true
stderr_logfile=/var/log/fpm.err.log
stdout_logfile=/var/log/fpm.out.log